# user/.config/systemd/user/timers.target.wants/heartbeat.timer
[Unit]
Description=Do heartbeat every minutes.

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
