# etc/root.crontab
30 0 * * * /sbin/reboot
30 * * * * /liims/bin/reset.sh 
