# etc/root.crontab
0 0 * * Sat /sbin/reboot
30 * * * * /liims/bin/reset.sh 
