# etc/systemd/system/liims-reset.service
[Unit]
Description=Reset LIIMS user

[Service]
Type=oneshot
ExecStart=/usr/local/bin/liims-reset.sh
