# user/.config/systemd/user/chameleon.service
[Unit]
Description=LIIMS chameleon: change midori settings by boot cmdline
Before=midori.service

[Service]
Type=oneshot
ExecStart=/liims/bin/chameleon.sh
RemainAfterExit=yes
