# user/.config/systemd/user/heartbeat.service
[Unit]
Description=LIIMS heartbeat.

[Service]
Type=oneshot
ExecStart=/liims/bin/heartbeat 
