# user/.config/systemd/user/heartbeat.timer
[Unit]
Description=Do heartbeat every minutes.

[Timer]
OnBootSec=30
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
